# etc/profile.d/go.sh
GOROOT=/usr/local/go
GOPATH=$HOME/go
export GOROOT
export GOPATH
PATH=$GOPATH/bin:$GOROOT/bin:$PATH
